# etc/profile.d/texlive.sh
export MANPATH=/media/storage/texlive/2020/texmf-dist/doc/man/:$MANPATH
export INFOPATH=/media/storage/texlive/2020/texmf-dist/doc/info/:$INFOPATH
export PATH=/media/storage/texlive/2020/bin/x86_64-linux/:$PATH
