# home/.config/plasma-backup/autostart/xset.desktop
[Desktop Entry]
Exec=xset r rate 400 50
Icon=
Name=xset
Path=
Terminal=False
Type=Application
